# eveabacus-frontend.service
[Unit]
Description=EVE Abacus Next.js Frontend
After=network.target

[Service]
WorkingDirectory=/var/www/eveabacus/frontend/standalone
ExecStartPre=/bin/bash -c "ls -la /var/www/eveabacus/frontend/standalone/ && echo BUILD_ID check: && ls -la /var/www/eveabacus/frontend/standalone/.next/BUILD_ID || echo BUILD_ID not found"
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
SyslogIdentifier=eveabacus-frontend
User=www-data
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target 